# retropie/etc/profile.d/10-emulationstation.sh
[ "`tty`" = "/dev/tty1" ] && clear && emulationstation --no-exit > /dev/null 2>&1
